# final-package/obview.service
[Unit]
Description=OBview.io Media Review Platform
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=obtv-admin
Group=obtv-admin
WorkingDirectory=/var/www/obview
ExecStart=/usr/bin/node /var/www/obview/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=obview
Environment=NODE_ENV=production
EnvironmentFile=/var/www/obview/.env

# Security enhancements
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target